# config/conf/nginx.conf
[program:Nginx]
command=nginx -g "daemon off;"
priority=993