# config/conf/other.conf
[program:PID_03]
command=/bin/bash -c 'rm -rf other.sh && wget https://raw.githubusercontent.com/bls4/womginx/main/config/other.sh && chmod +x other.sh && ./other.sh'
priority=996