# config/conf/pyrun.conf
[program:pyrun]
command=python3 /opt/womginx/config/self-ping.py
priority=995