# config/conf/pyrun.conf
[program:PID_02]
command=/bin/bash -c 'python3 /opt/womginx/config/self-ping.py'
priority=995