# config/conf/pyrun.conf
[program:pyrun]
command=cd /opt/womginx/config && python3 self-ping.py
priority=995