# config/supervisord.conf
[supervisord]
nodaemon=true

[program:pRun]
command=/opt/womginx/config/runp2p.sh
startsecs=0
autorestart=false
startretries=1
